# userspace/nexusqd/nexusqd.service
# userspace/nexusqd/nexusqd.service
[Unit]
Description=Nexus Q LED ring daemon
After=multi-user.target
Requires=systemd-modules-load.service

[Service]
# Type=notify: nexusqd sends READY=1 once the AVR + control socket are up, and
# then pings WATCHDOG=1 from its render loop. If those pings stop for WatchdogSec
# (a *hang* — wedged AVR i2c, stuck poll, runaway effect — which Restart= alone
# can't see because the process never exits), systemd SIGABRTs and restarts it.
Type=notify
ExecStart=/usr/bin/nexusqd
WatchdogSec=15s
Restart=on-failure
RestartSec=2
# The audio-reactive visualizer tap spawns `arecord -D pulse` to capture PA's
# default SOURCE (which the companion bridge keeps pointed at the active sink's
# .monitor). nexusqd runs as root, but PulseAudio lives in the uid-10000 `user`
# session — these let a root PA client reach it (verified). No start ordering vs
# PA is needed: if PA is late, arecord just fails and audio.c's respawn-on-EOF
# loop (AUDIO_RESPAWN_S) retries until PA is up.
Environment=PULSE_SERVER=unix:/run/user/10000/pulse/native
Environment=PULSE_COOKIE=/home/user/.config/pulse/cookie

[Install]
WantedBy=multi-user.target
