# pmos/device-google-steelhead/librespot.service
[Unit]
Description=librespot — Spotify Connect (Nexus Q)
Documentation=https://github.com/librespot-org/librespot
After=network-online.target sound.target
Wants=network-online.target

[Service]
# Output via the "nexusq" ALSA PCM (asound.conf): plug -> hw:1,0 forced to 48000 Hz.
# The McBSP2/TAS5713 link only clocks the 48 kHz family cleanly, so 44.1 kHz
# Spotify content is resampled to 48k; with the kernel CLKGDV fix FSYNC is an
# exact 48 kHz so playback is at correct pitch/tempo.
# --ap-port 443 works around VLAN20 blocking librespot's default AP port 4070.
# --zeroconf-port is fixed so the nftables drop-in (60_spotify.nft) can match it.
ExecStart=/usr/bin/librespot \
    --name "Nexus Q" \
    --backend alsa \
    --device nexusq \
    --bitrate 320 \
    --format S16 \
    --initial-volume 60 \
    --ap-port 443 \
    --zeroconf-port 37879 \
    --cache /var/cache/librespot
Restart=on-failure
RestartSec=5
RuntimeDirectory=librespot
CacheDirectory=librespot

[Install]
WantedBy=multi-user.target
