# pmos/device-google-steelhead/librespot.service
# librespot Spotify Connect — a systemd USER unit (runs in the uid-10000 `user`
# session so it shares that session's PulseAudio). Installed to
# /usr/lib/systemd/user/ and auto-started via a default.target.wants symlink
# (same mechanism as pulseaudio.service — no linger needed under autologin).
#
# This REPLACES the old root/system unit that went straight to ALSA
# (--device nexusq_soft, the asound.conf `type multi` fan-out) and shared its
# volume with the companion bridge via the "NexusQ" ALSA softvol. Under the
# PA-centric design librespot is just one PA INPUT; OUTPUT selection (TAS5713 /
# SPDIF / HDMI) and volume live in PulseAudio. All the args (PA device, the
# libmdns wlan0-interface pin, --ap-port 443, cache) are in the wrapper
# /usr/bin/librespot-nexusq.
[Unit]
Description=librespot — Spotify Connect (Nexus Q, PulseAudio input)
Documentation=https://github.com/librespot-org/librespot
After=pulseaudio.service network-online.target
Wants=pulseaudio.service network-online.target
# This user unit is enabled globally (/etc/systemd/user/default.target.wants),
# so it would also start in root's (or any admin's) user manager on login and
# crash-loop fighting the appliance instance for the zeroconf port. Pin it to
# the appliance user (uid 10000); other users' managers skip it silently.
ConditionUser=10000

[Service]
ExecStart=/usr/bin/librespot-nexusq
Restart=on-failure
RestartSec=5
CacheDirectory=librespot

[Install]
WantedBy=default.target
