# pmos/device-google-steelhead/nexusq-usb-gadget.service
[Unit]
Description=Nexus Q: USB gadget = RNDIS net (172.16.42.1 + ssh) + ACM console
Documentation=file:///usr/bin/nexusq-usb-gadget.sh
# Own the gadget ourselves rather than relying on pmOS's RNDIS bringup + the
# fragile RNDIS->ACM swap (which could leave it unbound). Run after modules are
# loaded; the script waits for the musb UDC itself.
After=systemd-modules-load.service local-fs.target
# This replaces the swap-to-ACM console service; do not let both fight over g1.
Conflicts=nexusq-usb-acm-console.service

[Service]
Type=oneshot
ExecStart=/usr/bin/nexusq-usb-gadget.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
