# pmos/device-google-steelhead/nexusq-usb-acm-console.service
[Unit]
Description=Nexus Q: USB gadget as ACM serial console (replaces RNDIS network)
Documentation=https://github.com/petronijus/nexusQ-reloaded
# pmOS brings the micro-USB gadget up as an RNDIS *network* device (usb0) late in
# boot, when the musb UDC is ready (~90 s on this 350 MHz board). Trigger off that
# net device so we run exactly when the gadget exists, not before (a plain
# ConditionPathExists on g1 races and is skipped). Each boot pmOS makes RNDIS
# first, then this swaps it to ACM.
After=sys-subsystem-net-devices-usb0.device

[Service]
Type=oneshot
ExecStart=/usr/bin/nexusq-usb-acm-console.sh
# Offer a login on the freshly created gadget serial port.
ExecStartPost=-/bin/systemctl start serial-getty@ttyGS0.service
RemainAfterExit=yes

[Install]
# Pulled in by the usb0 net device appearing (systemd-native, no udev rule file).
WantedBy=sys-subsystem-net-devices-usb0.device
