# pmos/device-google-steelhead/nq-healthd.service
[Unit]
Description=Nexus Q runtime health monitor (cpufreq/power/thermal/LED ring/nexusqd)
Documentation=file:///usr/bin/nq-healthd
# Start late; nexusqd/librespot are part of what we watch.
After=multi-user.target nexusqd.service

[Service]
Type=simple
ExecStart=/usr/bin/nq-healthd
# Self-heal: the monitor itself must never be the thing that dies silently.
Restart=always
RestartSec=5
# Stay out of the way on this single-/dual-core OMAP4.
Nice=10
IOSchedulingClass=idle
# Logs live under /var/log/nq-health (rotated by the daemon itself).

[Install]
WantedBy=multi-user.target
