# userspace/nexusq-setupd/nexusq-setupd.service
[Unit]
Description=Nexus Q BT WiFi-provisioning (setup mode)
# BT + NM must exist; none is a hard dep (setupd degrades gracefully).
# NO After=nexusqd.service: nexusqd is After=multi-user.target, so ordering on
# it from a WantedBy=multi-user.target unit creates a cycle and systemd DELETES
# this unit's start job (bit us in v1.9.0-rc1; see HANDOFF 2026-07-13).
After=bluetooth.service NetworkManager.service
Wants=bluetooth.service

[Service]
Type=simple
ExecCondition=/usr/bin/nexusq-setup-needed
ExecStart=/usr/bin/nexusq-setupd
# One shot per trigger: it exits after finishSetup / inactivity; a crash
# restarts it (setup must survive daemon bugs while the user is mid-flow).
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
